# templates/etc/profile.d/java.sh.j2
export JAVA_HOME={{java_home}}
export JRE_HOME={{java_home}}
